# overlays/volla-yggdrasil/etc/rc.local
#!/bin/sh

# lower panel backlight a bit more
echo 50 > /sys/class/leds/lcd-backlight/brightness || :
